# r170fancontroller.service
[Unit]
Description=R710fancontroller

[Service]
ExecStart=/usr/sbin/r710fancontroller.sh

[Install]
WantedBy=yourLoginHere
